# scheduler/cs-evolver.timer
[Unit]
Description=Trigger CyberSentinel Evolver weekly evolution loop

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
